# deploy/no-docker/tg-swarm.service
[Unit]
Description=tg-swarm FastAPI runtime
After=network.target

[Service]
Type=simple
User=tg-swarm
Group=tg-swarm
WorkingDirectory=/opt/tg-swarm/current
EnvironmentFile=/etc/tg-swarm.env
ExecStart=/opt/tg-swarm/.venv/bin/python /opt/tg-swarm/current/server.py
Restart=always
RestartSec=5
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=full
ReadWritePaths=/var/lib/tg-swarm

[Install]
WantedBy=multi-user.target
